# ktrade-backend.service
# KTrade dashboard backend (Flask). Serves the UI + /all + manual-order routes.
#   sudo cp deploy/ktrade-backend.service /etc/systemd/system/
#   sudo systemctl enable --now ktrade-backend
[Unit]
Description=KTrade backend (Flask dashboard API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ktrade
WorkingDirectory=/app
EnvironmentFile=/app/.env
ExecStart=/app/.venv/bin/python backend/ktrade_alpaca.py
Restart=always
RestartSec=10
StandardOutput=append:/app/ktrade_backend.log
StandardError=append:/app/ktrade_backend.log

[Install]
WantedBy=multi-user.target
